# crontab
@hourly /root/com-list.sh ; chgrp users /tmp/com-menu.txt

